# relay/sources/rudder-relay.cron
# Cron file for Rudder relay
#

*/5 * * * * root /opt/rudder/bin/relay-cleanup >/dev/null

