# relay/sources/rudder-relay.cron
# Cron file for Rudder relay
#

*/5 * * * * root /opt/rudder/share/relay-api/cleanup.sh >/dev/null

